# tokotanionline-nextjs/deploy/tokotani-engine-hub.service
# FASE 7.2 — SYSTEMD SERVICE: Go Engine Hub
# 
# Installation:
#   sudo cp deploy/tokotani-engine-hub.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tokotani-engine-hub
#   sudo systemctl start tokotani-engine-hub
#
# Management:
#   sudo systemctl status tokotani-engine-hub
#   sudo systemctl restart tokotani-engine-hub
#   sudo journalctl -u tokotani-engine-hub -f

[Unit]
Description=Tokotani Online Engine Hub (Go Backend)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/tokotani-online-nextjs/engine-hub

# FASE 7.1: Production ENV - All config from OS environment
# DO NOT use .env file in production
# Set environment variables in /etc/systemd/system/tokotani-engine-hub.service.d/override.conf
Environment="ENV=production"
Environment="ENGINE_PORT=8090"

# Required ENV vars (set via systemd environment or EnvironmentFile):
# - OPENAI_API_KEY (or AI_API_KEY) - CRITICAL, server will not start without this
# - DATABASE_URL (if using database features)

ExecStart=/opt/tokotani-online-nextjs/engine-hub/bin/server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tokotani-engine-hub

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=150%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tokotani-online-nextjs/engine-hub/logs

[Install]
WantedBy=multi-user.target
